# .docker/supervisor/supervisor-horizon.conf
[supervisord]
nodaemon=true

[program:horizon]
command=php /var/www/artisan queue:work
numprocs=1
autostart=true
autorestart=true
user=www-data